# config/profile.d/step_10_deploy.sh






cp /home/kas/kas_workspace/zhaoliming/GitProjects/dotfiles/config/profile.d/profile /etc/profile
cp /home/kas/kas_workspace/zhaoliming/GitProjects/dotfiles/config/profile.d/setup_env.sh /etc/profile.d/setup_env.sh

chmod 644 /etc/profile /etc/profile.d/setup_env.sh